# docs/wa/wa.cron

*/20 6-20 * * * /usr/bin/bash -x /home/admin/Documents/basic/docs/wa/deploy.sh & 
