# docs/wa/wa.cron

*/18,36,49 6-20 * * *	/usr/bin/bash -x /home/admin/Documents/basic/docs/wa/deploy.sh & 
